# deploy/devbloom-web.service.example
# Example systemd unit for DevBloom web (Next.js).
# Copy to /etc/systemd/system/devbloom-web.service
# Replace APP_ROOT paths below if your clone lives elsewhere; set USER/Group if not ec2-user.
# Set NEXT_PUBLIC_API_URL_BASE to your API URL (path-based: https://yourdomain.com/api, or subdomain: https://api.yourdomain.com).
# Then: sudo systemctl daemon-reload && sudo systemctl enable --now devbloom-web

[Unit]
Description=DevBloom (Next.js)
After=network.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/github/devbloom/web
Environment=NODE_ENV=production
Environment=NEXT_PUBLIC_API_URL_BASE=https://dev.funbloomstudio.com/api
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
