# deploy/gamedev-web.service.example
# Example systemd unit for Game Dev King web (Next.js).
# Copy to /etc/systemd/system/gamedev-web.service
# Replace /home/ec2-user/gamedev-king with your APP_ROOT if different; set USER/Group if not ec2-user.
# Set NEXT_PUBLIC_API_URL_BASE to your API URL (path-based: https://yourdomain.com/api, or subdomain: https://api.yourdomain.com).
# Then: sudo systemctl daemon-reload && sudo systemctl enable --now gamedev-web

[Unit]
Description=Game Dev King Web (Next.js)
After=network.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/gamedev-king/web
Environment=NODE_ENV=production
Environment=NEXT_PUBLIC_API_URL_BASE=https://yourdomain.com/api
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
