# deploy/gamedev-api.service.example
# Example systemd unit for Game Dev King API (FastAPI).
# Copy to /etc/systemd/system/gamedev-api.service
# Replace APP_ROOT with your app directory (e.g. /home/ec2-user/gamedev-king) and USER if not ec2-user.
# Then: sudo systemctl daemon-reload && sudo systemctl enable --now gamedev-api

[Unit]
Description=DevBloom API
After=network.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/gamedev-king/api
EnvironmentFile=/home/ec2-user/gamedev-king/api/.env
ExecStart=/home/ec2-user/gamedev-king/api/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
